# backend/deploy/aipic.service
[Unit]
Description=AI Image Generation Service (aipic)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/aipic/backend
EnvironmentFile=/opt/aipic/backend/.env
ExecStart=/opt/aipic/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8800 --workers 2
Restart=always
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aipic/backend/outputs /opt/aipic/backend/temp /opt/aipic/backend/user_data /opt/aipic/backend/logs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=aipic

[Install]
WantedBy=multi-user.target
